# resolvrServer.conf
[program:resolvr-server]
command=nodejs /srv/resolvr/server.js
directory=/srv/resolvr
numprocs=1
user=root
process_name=resolvrServer
autostart=true
autorestart=true
redirect_stderr=true


